# libs/assistant/templates/homeassistant.service.erb
[Unit]
Description=Home Assistant
After=network-online.target

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=<%= @config_dir %>
Environment="PATH=<%= @venv_dir %>/bin:/usr/local/bin:/usr/bin:/usr/local/bin/uv"
ExecStart=<%= @venv_dir %>/bin/python3 -m homeassistant --config <%= @config_dir %>
Restart=always
RestartForceExitStatus=100

[Install]
WantedBy=multi-user.target
