# templates/runner.service.erb
[Unit]
Description=Gitea Actions Runner
After=network.target

[Service]
Type=simple
User=<%= node['git']['app']['user'] %>
Group=<%= node['git']['app']['group'] %>
WorkingDirectory=<%= node['runner']['install_dir'] %>
ExecStart=<%= node['runner']['install_dir'] %>/ace_runner daemon --config <%= node['runner']['install_dir'] %>/config.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
